# nginx-proxy/supervisor/confd.conf
[program:confd]
command=/bin/sh -c "usr/bin/confd -backend etcd -node $ETCD_HOST"
