# nginx-proxy/supervisor/confd.conf
[program:confd]
command=/bin/bash -c "/usr/bin/confd -interval=5 -node=http://${ETCD_HOST}"
